# gamecenter/etc/gamecenter_api.ini
[group:warriors]
programs=gamecenter-8004, gamecenter-8005, gamecenter-8006, gamecenter-8007

[program:gamecenter-8004]
command=/home/game_venv/bin/game-api --port=8004
directory=/root/gamecenter
autostart=true
autorestart=true
startretries=999
redirect_stderr=true
stdout_logfile=/var/log/supervisor/game-api.log
stderr_logfile=/var/log/supervisor/game-api-err.log
stdout_logfile_maxbytes = 20MB
stderr_logfile_maxbytes = 20MB

[program:gamecenter-8005]
command=/home/game_venv/bin/game-api --port=8005
directory=/root/gamecenter
autostart=true
autorestart=true
startretries=999
redirect_stderr=true
stdout_logfile=/var/log/supervisor/game-api.log
stderr_logfile=/var/log/supervisor/game-api-err.log
stdout_logfile_maxbytes = 20MB
stderr_logfile_maxbytes = 20MB

[program:gamecenter-8006]
command=/home/game_venv/bin/game-api --port=8006
directory=/root/gamecenter
autostart=true
autorestart=true
startretries=999
redirect_stderr=true
stdout_logfile=/var/log/supervisor/game-api.log
stderr_logfile=/var/log/supervisor/game-api-err.log
stdout_logfile_maxbytes = 20MB
stderr_logfile_maxbytes = 20MB

[program:gamecenter-8007]
command=/home/game_venv/bin/game-api --port=8007
directory=/root/gamecenter
autostart=true
autorestart=true
startretries=999
redirect_stderr=true
stdout_logfile=/var/log/supervisor/game-api.log
stderr_logfile=/var/log/supervisor/game-api-err.log
stdout_logfile_maxbytes = 20MB
stderr_logfile_maxbytes = 20MB