# etc/gamecenter_worker.ini
[program:gamecenter-worker]
command=/home/game_venv/bin/game-worker
directory=/home/gamecenter
autostart=true
autorestart=true
startretries=999
redirect_stderr=true
stdout_logfile=/var/log/supervisor/game-worker.log
stderr_logfile=/var/log/supervisor/game-worker-err.log
stdout_logfile_maxbytes = 20MB
stderr_logfile_maxbytes = 20MB
